# docs/deploy/server4-codex-buddy-cron.timer
[Unit]
Description=Run codex-buddy cleanup and reminders

[Timer]
OnBootSec=3min
OnUnitActiveSec=10min
RandomizedDelaySec=30s
Persistent=true

[Install]
WantedBy=timers.target
