# docs/deploy/server4-codex-buddy-cron.service
[Unit]
Description=codex-buddy cleanup and reminder pass

[Service]
Type=oneshot
WorkingDirectory=/opt/projects/codex-hud
EnvironmentFile=/etc/codex-buddy/codex-buddy.env
ExecStart=/usr/bin/npm --workspace @codex-hud/buddy-server run cron
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=read-only
ReadWritePaths=/var/lib/codex-buddy
